# jack-coreutils/src/systemd/jack-coreutils.service
[Unit]
Description=Lichee-Jack coreutils daemon

[Service]
Type=simple
ExecStart=/usr/bin/jackstart

[Install]
WantedBy=multi-user.target

